# supervisord.conf
[supervisord]
nodaemon=true

[program:shellinabox]
command=/bin/bash -c  /scripts/shellinabox.sh

[program:lighttpd]
command = /usr/local/sbin/lighttpd -D
